# .config/autostart/bras.desktop
[Desktop Entry]
Exec=braslogin 9401925 83596085
Hidden=true
Name=bras
Type=Application
Version=1.0
X-Razor-Need-Tray=true
