# .config/autostart/compiz.desktop
[Desktop Entry]
Exec=compiz --replace ccp
Name=compiz
Type=Application
Version=1.0
