# .config/autostart/compiz.desktop
[Desktop Entry]
Exec=compiz --replace ccp
Hidden=true
Name=compiz
Type=Application
Version=1.0
