# .config/autostart/desktop.desktop
[Desktop Entry]
Exec=pcmanfm --desktop
Name=desktop
Type=Application
Version=1.0
