# .config/autostart/eHosts.desktop
[Desktop Entry]
Exec=eHosts -s
Name=eHosts
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
