# .config/autostart/eHosts.desktop
[Desktop Entry]
Exec=eHosts -s -f /etc/hosts -d 223.5.5.5
Name=eHosts
Type=Application
Version=1.0
