# .config/autostart/Google music.desktop
[Desktop Entry]
Exec=google-musicmanager 
Name=Google music
Type=Application
Version=1.0
X-Razor-Need-Tray=true
