# .config/autostart/sshR.desktop
[Desktop Entry]
Exec=sshr
Name=sshr
Type=Application
Version=1.0
X-Razor-Need-Tray=true
