# .config/autostart/synapse.desktop
[Desktop Entry]
Encoding=UTF-8
Exec=synapse --startup
Icon=synapse
Name=Synapse
Type=Application
X-GNOME-Autostart-enabled=true
