# .config/autostart/Synapse.desktop
[Desktop Entry]
Exec=synapse -s
Name=Synapse
Type=Application
Version=1.0
