# .config/autostart/tilda.desktop
[Desktop Entry]
Exec=tilda
Name=tilda
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
