# .config/autostart/网络链接.desktop
[Desktop Entry]
Exec=nm-applet
Name=网络链接
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
